# src/mirror/mirrors.crontab
15 2,5,8,11,14,17,20,23 * * *   nornetpp   flock -x -n /var/lock/nornet-mirror-f18-1.lock -c "/etc/nornet/fedora-mirror.sh -f f18 -a x86_64 -e -r /filesrv/adm/mirror/ -l"
15 1,4,7,10,13,16,19,22 * * *   nornetpp   flock -x -n /var/lock/nornet-mirror-f18-2.lock -c "/etc/nornet/fedora-mirror.sh -f f18 -a i386   -e -r /filesrv/adm/mirror/ -l"
45 2,5,8,11,14,17,20,23 * * *   nornetpp   flock -x -n /var/lock/nornet-mirror-f20-1.lock -c "/etc/nornet/fedora-mirror.sh -f f20 -a x86_64 -e -r /filesrv/adm/mirror/ -l"
45 1,4,7,10,13,16,19,22 * * *   nornetpp   flock -x -n /var/lock/nornet-mirror-f20-2.lock -c "/etc/nornet/fedora-mirror.sh -f f20 -a i386   -e -r /filesrv/adm/mirror/ -l"
