# src/mirror/mirrors.crontab
# Fedora Core 18
# Alternative: rsync://ftp.uninett.no/fedora-linux-releases/18/Everything/
30 6,21 * * *                   mkdir -p /nfs/adm/mirror/fedora/releases/18/Everything/ && rsync --archive --compress --delete --delete-excluded --no-motd --exclude iso/ rsync://ftp.uninett.no/fedora-linux-releases/18/Everything/ /nfs/adm/mirror/fedora/releases/18/Everything/
45 0,3,6,9,12,15,18,21 * * *    mkdir -p /nfs/adm/mirror/fedora/releases/18/Everything/ && rsync --archive --compress --delete --delete-excluded --no-motd --exclude iso/ rsync://ftp.uninett.no/fedora-linux-updates/18/ /nfs/adm/mirror/fedora/updates/18/

# Ubuntu 12.04
# Alternative: rsync://ftp.uninett.no/ubuntu/
30 1,4,7,10,13,16,19,22 * * * mkdir -p /nfs/adm/mirror/ubuntu/ && rsync --archive --compress --delete --delete-excluded --no-motd rsync://se.archive.ubuntu.com/ubuntu/ /nfs/adm/mirror/ubuntu/
