# src/Mirror/nornet-mirror.crontab
15 2,5,8,11,14,17,20,23 * * *   nornetpp   flock -x -n /var/lock/nornet-mirror-f18-1.lock -c "/etc/nornet/fedora-mirror.sh -u rsync://archives.fedoraproject.org/fedora-archive/fedora/linux -f f18 -a x86_64 -r /filesrv/adm/mirror/ -l"
15 1,4,7,10,13,16,19,22 * * *   nornetpp   flock -x -n /var/lock/nornet-mirror-f18-2.lock -c "/etc/nornet/fedora-mirror.sh -u rsync://archives.fedoraproject.org/fedora-archive/fedora/linux -f f18 -a i386   -r /filesrv/adm/mirror/ -l"
# 30 2,5,8,11,14,17,20,23 * * *   nornetpp   flock -x -n /var/lock/nornet-mirror-f20-1.lock -c "/etc/nornet/fedora-mirror.sh -e -f f20 -a x86_64 -r /filesrv/adm/mirror/ -l"
# 30 1,4,7,10,13,16,19,22 * * *   nornetpp   flock -x -n /var/lock/nornet-mirror-f20-2.lock -c "/etc/nornet/fedora-mirror.sh -e -f f20 -a i386   -r /filesrv/adm/mirror/ -l"
45 2,5,8,11,14,17,21,23 * * *   nornetpp   flock -x -n /var/lock/nornet-mirror-f21-1.lock -c "/etc/nornet/fedora-mirror.sh -e -f f21 -a x86_64 -r /filesrv/adm/mirror/ -l"
45 1,4,7,10,13,16,19,22 * * *   nornetpp   flock -x -n /var/lock/nornet-mirror-f21-2.lock -c "/etc/nornet/fedora-mirror.sh -e -f f21 -a i386   -r /filesrv/adm/mirror/ -l"
00 2,5,8,11,14,17,21,23 * * *   nornetpp   flock -x -n /var/lock/nornet-mirror-f22-1.lock -c "/etc/nornet/fedora-mirror.sh -e -f f22 -a x86_64 -r /filesrv/adm/mirror/ -l"
00 1,4,7,10,13,16,19,22 * * *   nornetpp   flock -x -n /var/lock/nornet-mirror-f22-2.lock -c "/etc/nornet/fedora-mirror.sh -e -f f22 -a i386   -r /filesrv/adm/mirror/ -l"
