# src/tests/C2/rc.local
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

ip -6 addr add fd55:5555::b/64 dev eth3 || true
ip addr add 169.254.100.11/24 dev eth1 || true
ip addr add 172.31.255.11/24  dev eth1 || true
ip addr add 132.252.156.11/24 dev eth1 || true
ip route add 10.0.0.0/8     via 132.252.156.70 || true
ip route add 192.168.0.0/16 via 132.252.156.70 || true

touch /var/lock/subsys/local
