# debian/nornet-development.cron.d
45 5 * * *	root	/usr/bin/Random-Sleep 0 3600 -quiet && pbuilder update || pbuilder create
45 6 * * *	root	/usr/bin/Random-Sleep 0 3600 -quiet && apt-file update
