# debian/nornet-tunnelbox.cron.d
8,38 * * * *  root    test -x /etc/init.d/nornet-tunnelbox && /usr/bin/Random-Sleep 0 1200 -quiet && service nornet-tunnelbox check-and-configure
