# debian/nornet-tunnelbox.cron.d
8,28,48 * * * *  root    test -x /etc/init.d/nornet-tunnelbox && /usr/bin/Random-Sleep 0 150 -quiet && service nornet-tunnelbox check-and-configure
