# debian/nornet-development.cron.d
45 5 * * *	root	pbuilder update || sudo pbuilder create
45 6 * * *	root	apt-file update
