# debian/nornet-filesrv.cron.d
0,15,30,45 * * * *  root    test -x /etc/init.d/nornet-filesrv && /etc/init.d/nornet-filesrv check-and-configure
