# debian/nornet-monitor.cron.d
2,7,12,17,22,27,22,27,42,47,52,57 * * * *  root    test -x /etc/init.d/nornet-monitor && service nornet-monitor check-and-configure
