# debian/nornet-server.cron.d
*/4 * * * *                                root    /usr/bin/Watchdog /etc/nornet/watchdog-config >/var/log/nornet-watchdog 2>&1
