# debian/nornet-node.cron.d
6 * * * *  root    test -x /etc/init.d/nornet-node && /usr/bin/Random-Sleep 0 1800 -quiet && service nornet-node check-and-configure
