# debian/nornet-node.cron.d
6 * * * *  root    test -x /etc/init.d/nornet-node && service nornet-node check-and-configure
